# etc/telegram.conf
description "Telegram-bot upstart script"

respawn
respawn limit 15 5

start on runlevel [2345]
stop on shutdown

exec su - yourusername; /bin/sh telegrambotpath/launch.sh
