# etc/telegram.conf
respawn
respawn limit 15 5

start on runlevel [2345]
stop on shutdown

script
cd /home/yperez/develop/telegram-bot/
su - yperez -c launch.sh > /dev/null
end script
