# etc/telegram.conf
respawn
respawn limit 15 5

start on runlevel [2345]
stop on shutdown

script
su - yperez -c "screen /home/yperez/develop/telegram-bot/launch.sh"
end script
