# etc/telegram.conf
respawn
respawn limit 15 5

start on runlevel [2345]
stop on shutdown

script
exec su -s /bin/sh -c /home/yperez/develop/telegram-bot/launch.sh yperez
end script
