# etc/telegram.conf
respawn
respawn limit 15 5

start on runlevel [2345]
stop on shutdown

script
su - yperez -c "cd /home/yperez/develop/telegram-bot/ && screen ./launch.sh"
end script
